# DaaS-project/ansible/roles/desktop-config/templates/inactivity-monitor.service.j2
[Unit]
Description=Monitoring d'inactivité pour desktop ephemeral
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/inactivity-monitor.sh
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
